# lixyswarm-seed.service
[Unit]
Description=LixySwarm LSP v3 Seed
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lixyswarm
Group=lixyswarm
WorkingDirectory=/opt/lixyswarm/app
Environment=PYTHONPATH=/opt/lixyswarm/app
Environment=LIXYSWARM_IDENTITY_PATH=/var/lib/lixyswarm/identity.key
Environment=LIXYSWARM_PEERS_PATH=/var/lib/lixyswarm/peers_v3.json
EnvironmentFile=-/etc/lixyswarm/seed.env
ExecStart=/opt/lixyswarm/venv/bin/python /opt/lixyswarm/app/node_daemon.py
Restart=always
RestartSec=5
StateDirectory=lixyswarm
StateDirectoryMode=0700
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_INET AF_INET6
LockPersonality=true

[Install]
WantedBy=multi-user.target
